# scripts/power-monitor.base.service
[Unit]
Description=Power Monitor
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
Restart=always
RestartSec=10
User=